# scripts/program-backup.service
[Unit]
Description=BSA workout-program safety-net backup (versioned, 90d retention)

[Service]
Type=oneshot
User=ec2-user
Environment=HOME=/home/ec2-user
ExecStart=/home/ec2-user/run_backup.sh
